# src/isl/inittab
#
#	@(#)inittab	3.3	ULTRIX OSF	7/2/91
#
# (c) Copyright 1990, OPEN SOFTWARE FOUNDATION, INC.
# ALL RIGHTS RESERVED
#
#
# OSF/1 Release 1.0
#
#	001	02-jul-1991	ccb
#		this one allows autoboot to return to installation
#		interrupts are enabled in .profile (stty sane)
is:3:initdefault:
inst:23:wait:/sbin/init S < /dev/console > /dev/console 2>&1

